# php_env_https/supervisord.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
user=root
numprocs=1
stdout_logfile=/var/log/nginx/access.log